# data/me.blq.FedoraUpdater.check.service
[Unit]
Description=Check for Fedora system updates

[Service]
Type=oneshot
ExecStart=fedora-updater --check-updates
