# data/me.blq.FedoraUpdater.timer
[Unit]
Description=Periodic Fedora update check

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
